# public/data/winreg/artifacts/autostart-and-persistence-wmi-event-subscriptions-registry.json
{"slug":"autostart-and-persistence-wmi-event-subscriptions-registry","name":"WMI Event Subscriptions (Registry)","category":"autostart_and_persistence","categoryLabel":"Autostart & Persistence","categoryDescription":"Registry keys that cause code to execute at boot, logon, or via Windows subsystem hooks.","hive":["SOFTWARE"],"keys":["Microsoft\\WBEM\\ESS"],"description":"Registry-side WMI Event Subscription data — the WBEM\\ESS key reflects active permanent WMI subscriptions (EventFilter, EventConsumer, FilterToConsumerBinding).","forensic_value":"WMI subscriptions are one of the most persistent and detection-resistant mechanisms available. The registry key shows subscription names but the full definition is in the WMI repository (Objects.data). Cross-reference with WMI-Activity Operational log EID 5861.","mitre":"TA0003","techniques":["T1546.003"],"parsers":["Autoruns","RegRipper"],"source":"github.com/dfir-scripts/dfir-scripts.github.io","sourceUrl":"https://dfir-scripts.github.io/registry/","license":"MIT"}