# public/data/winreg/artifacts/autostart-and-persistence-active-setup-installed-components.json
{"slug":"autostart-and-persistence-active-setup-installed-components","name":"Active Setup — Installed Components","category":"autostart_and_persistence","categoryLabel":"Autostart & Persistence","categoryDescription":"Registry keys that cause code to execute at boot, logon, or via Windows subsystem hooks.","hive":["SOFTWARE"],"keys":["Microsoft\\Active Setup\\Installed Components"],"description":"Active Setup entries run once per user on first logon before the shell loads. StubPath command executes if per-user version counter is lower than system version.","forensic_value":"Classic persistence — adding a subkey causes StubPath to run as the user on next logon without Run-key detection.","mitre":"TA0003","techniques":["T1547.014"],"parsers":["RegRipper (installedcomp.pl)","Autoruns"],"source":"github.com/dfir-scripts/dfir-scripts.github.io","sourceUrl":"https://dfir-scripts.github.io/registry/","license":"MIT"}