# public/data/daily-briefs/cyber/2026-07-23.json
{
  "type": "cyber",
  "date": "2026-07-23",
  "threatLevel": "Rockwell Automation / Schneider Electric / Siemens PLCs",
  "executiveSummary": "Active OT targeting by Iranian‑affiliated actors is exploiting internet‑exposed PLCs (Rockwell, Schneider, Siemens), causing operational disruption across Water, Energy, and Government Services. Remove direct internet access, validate controller logic integrity, and enforce access controls immediately. Security Advisory Rapid exploitation of on‑prem Microsoft SharePoint deserialization flaws enables machine‑key theft and RCE, with CISA urging hardening and patching; multiple SharePoint and ADFS items are now in KEV. Expect lateral movement toward OT via identity compromise. Security Advisory KEV Update Authentication bypasses on critical VPN gateways (Check Point Remote Access; PAN‑OS GlobalProtect CAS profiles) risk direct ingress to enterprise and potentially OT networks. Disable legacy/profiles and patch urgently. Alert Alert EV charging infrastructure (EVbee DC‑80) has multiple unauthenticated RCE/command injection paths (web, Bluetooth, OCPP), creating safety and grid‑impact risk; immediate isolation and hardening are required. View Article Nation‑state and proxy infrastructure: China‑nexus actors (incl. Volt Typhoon) are leveraging covert SOHO/IoT proxy networks for CNI pre‑positioning; allied advisory urges urgent edge hardening and egress controls. Supporting Report CERT Polska disclosed critical flaws in ICU Scandinavia Boomerang monitoring: plaintext credentials via HTTP and unauthorized sensor DB writes—patch and rotate creds immediately. View Article Multiple ICS‑relevant CVEs (Delta PLC Modbus, PcVue credentials, DAQFactory RCE, EIP/CIP stack DoS) raise near‑term risk to manufacturing and water sector operations. CVE Summary CVE Summary Fortinet perimeter exposure: large‑scale FortiGate credential leakage (“FortiBleed”) with continued validity in some cases—rotate secrets, review logs, and enforce PBKDF2 hashing. Alert",
  "keyFindings": [],
  "dashboard": { "kpis": [], "activelyExploited": [], "vendors": [], "sectors": [] },
  "topThreats": [],
  "threatActors": [],
  "cveWatch": [],
  "events": [],
  "ttps": {
    "descriptions": [
      "Initial Access: VPN auth bypass (T1133 External Remote Services), Exploit Public‑Facing App (T1190), Exposed VNC/HMI (T1021.001)",
      "Execution/Persistence: Webshells on SharePoint/GeoServer (T1505.003), Scheduled Tasks (T1053), DLL search‑order hijacking (Windows) (T1574.001)",
      "Defense Evasion: OAuth/Service Principal abuse in M365 Graph (T1550.001/T1098.001), AMSI/WAF evasion via deserialization",
      "Command & Control: DNS tunneling/AAAA channeling (T1071.004), Covert SOHO/IoT proxy chains",
      "ICS‑Specific: Unauthorized Modbus/CIP commands, PLC mode changes, AOI tampering; OCPP DataTransfer injection"
    ],
    "mitreIds": [
      "T1133",
      "T1190",
      "T1021.001",
      "T1505.003",
      "T1053",
      "T1574.001",
      "T1550.001",
      "T1098.001",
      "T1071.004"
    ]
  },
  "outlook72h": "Expect continued automated probing of SharePoint endpoints following PoC circulation; prioritize patch/validation and WAF rulesets on exposed portals. Monitor GlobalProtect and Check Point RA auth logs for anomalous sessions; temporarily restrict access to allow‑listed IPs where feasible. Increase OT perimeter monitoring for Modbus/CIP anomalies and VNC scans; enforce account lockouts and credential rotation on internet‑facing services. EV charging/OCPP operators should deploy emergency ACLs and mTLS enforcement; review firmware signing policies.",
  "relatedCves": ["CVE-2026-46458", "CVE-2026-50751", "CVE-2026-0265", "CVE-2024-36401", "CVE-2026-8933"],
  "rawMarkdown": "OT/ICS Cyber Threat Intelligence Brief OT/ICS Cyber Threat Intelligence Brief Report date: 2026-07-23 Overall OT cyber threat level: HIGH 17 C:ESCALATE events Immediate action 28 B:MONITOR events Heightened vigilance 10 A:IGNORE or low-operational-value items No action Executive Summary Active OT targeting by Iranian‑affiliated actors is exploiting internet‑exposed PLCs (Rockwell, Schneider, Siemens), causing operational disruption across Water, Energy, and Government Services. Remove direct internet access, validate controller logic integrity, and enforce access controls immediately. Security Advisory Rapid exploitation of on‑prem Microsoft SharePoint deserialization flaws enables machine‑key theft and RCE, with CISA urging hardening and patching; multiple SharePoint and ADFS items are now in KEV. Expect lateral movement toward OT via identity compromise. Security Advisory KEV Update Authentication bypasses on critical VPN gateways (Check Point Remote Access; PAN‑OS GlobalProtect CAS profiles) risk direct ingress to enterprise and potentially OT networks. Disable legacy/profiles and patch urgently. Alert Alert EV charging infrastructure (EVbee DC‑80) has multiple unauthenticated RCE/command injection paths (web, Bluetooth, OCPP), creating safety and grid‑impact risk; immediate isolation and hardening are required. View Article Nation‑state and proxy infrastructure: China‑nexus actors (incl. Volt Typhoon) are leveraging covert SOHO/IoT proxy networks for CNI pre‑positioning; allied advisory urges urgent edge hardening and egress controls. Supporting Report CERT Polska disclosed critical flaws in ICU Scandinavia Boomerang monitoring: plaintext credentials via HTTP and unauthorized sensor DB writes—patch and rotate creds immediately. View Article Multiple ICS‑relevant CVEs (Delta PLC Modbus, PcVue credentials, DAQFactory RCE, EIP/CIP stack DoS) raise near‑term risk to manufacturing and water sector operations. CVE Summary CVE Summary Fortinet perimeter exposure: large‑scale FortiGate credential leakage (“FortiBleed”) with continued validity in some cases—rotate secrets, review logs, and enforce PBKDF2 hashing. Alert Key Findings OT‑native threat: Direct PLC tampering and HMI data manipulation confirmed across U.S. critical sectors; actors abused native engineering tools and exposed services. Security Advisory Perimeter auth bypasses on enterprise VPN/portals (Check Point, PAN‑OS) create low‑friction paths to OT via IT gateways and jump hosts. View Article View Article EV infrastructure is exploitable at multiple layers (OCPP, firmware, Bluetooth), enabling unsafe charging control and charger outages. CVE Bulletin SharePoint on‑prem RCE is being exploited in the wild to steal IIS machine keys and persist via deserialization paths—hardening and credential rotation are essential. Security Advisory Facility monitoring software (ICU Boomerang) exposes plaintext credentials and unauth write access to sensor databases, risking plant data integrity. View Article Threat Dashboard 5 Actively exploited KEV items affecting edge/identity KEV focus 9 OT/ICS vendor advisories of note Patch planning 3 Nation‑state activity clusters (recent) Watchlisted Top Five Priority Threats 1) Iranian‑affiliated actors exploiting internet‑exposed PLCs (Rockwell, Schneider, Siemens) — Immediate removal from internet, logic integrity validation, and strict access controls are required to prevent further tampering and service disruption. Security Advisory 2) Check Point RA/Mobile Access CVE‑2026‑50751 (active exploitation) — Disable IKEv1 and legacy clients, require machine certificates, patch to fixed releases, and hunt for suspicious RA logins. Alert 3) PAN‑OS GlobalProtect CAS auth bypass CVE‑2026‑0265 — Remove CAS auth profiles from login interfaces until patched; restrict public exposure and enforce device certificates/MFA. Alert 4) SharePoint on‑prem RCE/deserialization (multiple CVEs; active exploitation) — Patch immediately, enable AMSI Full Mode, rotate machine keys/service creds, hunt for webshells. Security Advisory 5) EV charging infrastructure (EVbee DC‑80) multi‑vector compromise — Isolate chargers, enforce mutual TLS/OCPP controls, disable unsigned firmware updates, and monitor for OCPP DataTransfer abuse. CVE Bulletin OT/ICS Relevance Assessment Direct ICS impact confirmed (AA26‑097a): controller logic edits, AOI tampering, and HMI manipulation in multiple sectors. Perimeter to OT pivot risk remains elevated: active exploitation of SharePoint/ADFS and VPN auth bypasses can bridge to industrial DMZs and jump hosts. ICS protocol exposure persists (Modbus/CIP/S7): new PLC/stack CVEs and misconfigurations raise DoS and unauthorized command risk. Supply chain and remote‑service dependencies (EV charging, RUGGEDCOM APE + Virtual NGFW) can undermine OT perimeters if not segmented and hardened. Threat Actor Activity Iran‑affiliated groups targeting U.S. OT via exposed PLCs; use of leased infra and native vendor tools. Security Advisory China‑nexus (Volt Typhoon et al.) leveraging SOHO/IoT proxy fleets for CNI pre‑positioning. Advisory OilRig/APT34 “Project CAV3RN” abusing Microsoft Graph calendar events and DNS AAAA for C2—stealthy persistence in M365 estates. Research Pro‑Russia hacktivists exploiting exposed VNC to deface/alter HMI setpoints (loss of view). Advisory State exploitation of network edge/routers (PRC; FSB Center 16) to exfiltrate configs, harvest credentials, and mirror traffic. CSA CSA Vulnerability and CVE Watch Delta DVP12SE PLC Modbus unauth access/DoS; DTMSoft unsafe deserialization (RCE). CVE Summary PcVue ( CVE Summary EVbee DC‑80 multiple CVEs (web/Bluetooth/OCPP/firmware) enabling RCE and operational control. CVE Summary PAN‑OS CVE‑2026‑0265 (CAS) and Check Point CVE‑2026‑50751 (RA) — active/in‑scope auth bypasses. PAN‑OS Alert Check Point Alert SharePoint/ADFS multiple CVEs under active exploitation; machine‑key theft and RCE observed. CISA Tycon TPDIN‑Monitor‑WEB2 auth bypass/cleartext secrets; Siemens SIDIS SmartPlug critical vulns. ICS Advisory ICS Advisory ICU Scandinavia Boomerang (≤2.4.18.028) plaintext credentials and unauthorized writes. CERT Advisory Affected Vendors and Technologies Rockwell Automation / Schneider Electric / Siemens PLCs Palo Alto Networks PAN‑OS (GlobalProtect CAS) Check Point Security Gateways (RA/Mobile Access) Microsoft SharePoint / ADFS / IIS ICU Scandinavia Boomerang EVbee DC‑80 (OCPP/firmware/Bluetooth) Delta Electronics PLC & Tools PcVue HMI/SCADA Siemens RUGGEDCOM APE + Virtual NGFW Siemens SIDIS Secured SmartPlug Tycon TPDIN‑Monitor‑WEB2 GeoServer (CVE‑2024‑36401) Fortinet FortiGate / FortiSandbox Affected Sectors Water and Wastewater Systems; Energy (generation, T&D); Government Services and Facilities Transportation (EV charging infrastructure, telecom backbones); Critical Manufacturing Building Automation (KNX), Biometric/Facility systems Public sector GIS (GeoServer) supporting utilities and emergency services TTPs and MITRE ATT&CK-style observations Initial Access: VPN auth bypass (T1133 External Remote Services), Exploit Public‑Facing App (T1190), Exposed VNC/HMI (T1021.001) Execution/Persistence: Webshells on SharePoint/GeoServer (T1505.003), Scheduled Tasks (T1053), DLL search‑order hijacking (Windows) (T1574.001) Defense Evasion: OAuth/Service Principal abuse in M365 Graph (T1550.001/T1098.001), AMSI/WAF evasion via deserialization Command & Control: DNS tunneling/AAAA channeling (T1071.004), Covert SOHO/IoT proxy chains ICS‑Specific: Unauthorized Modbus/CIP commands, PLC mode changes, AOI tampering; OCPP DataTransfer injection Defensive Mitigation Priorities Remove OT devices from direct internet exposure; enforce industrial DMZ segmentation and monitored jump hosts. Security Advisory Patch and harden SharePoint/ADFS; enable AMSI Full Mode; rotate machine keys and service credentials; hunt for webshells. Security Advisory Disable legacy RA settings and patch VPN gateways (Check Point IKEv1; PAN‑OS CAS). Gate remote access with phishing‑resistant MFA. JPCERT JPCERT Isolate and harden EV charging/OCPP networks; enforce mTLS, disable unsigned firmware, and monitor for DataTransfer abuse. CISA Bulletin Rotate Fortinet/FortiGate credentials; upgrade to PBKDF2 hashing; audit for config exports and anomalous admin access. JPCERT Apply ICS advisories for Tycon TPDIN‑Monitor‑WEB2 and Siemens SIDIS SmartPlug; remove default accounts and minimize exposure. ICS Advisory ICS Advisory Patch ICU Boomerang to ≥2.4.18.029; rotate all exposed service/SMTP creds; restrict static file exposure. CERT Advisory Next 72‑Hour Outlook Expect continued automated probing of SharePoint endpoints following PoC circulation; prioritize patch/validation and WAF rulesets on exposed portals. Monitor GlobalProtect and Check Point RA auth logs for anomalous sessions; temporarily restrict access to allow‑listed IPs where feasible. Increase OT perimeter monitoring for Modbus/CIP anomalies and VNC scans; enforce account lockouts and credential rotation on internet‑facing services. EV charging/OCPP operators should deploy emergency ACLs and mTLS enforcement; review firmware signing policies. Event Cards C:ESCALATE C:ESCALATE Iranian‑Affiliated Actors Exploiting PLCs Across U.S. CNI THREAT: 5 Actors tampered with PLC logic and HMI data, causing operational impacts. Immediate actions: remove internet exposure, validate logic integrity, enforce change management, and segment OT. Rockwell / Schneider / Siemens Ports 44818/2222/102/502 Security Advisory C:ESCALATE SharePoint RCE & Machine‑Key Theft (Active) THREAT: 4 Active exploitation enables deserialization‑based RCE and key theft with persistent access and lateral movement potential. Patch, enable AMSI Full Mode, hunt for webshells, and rotate secrets. Microsoft SharePoint Webshells AMSI Security Advisory C:ESCALATE Check Point RA/Mobile Access CVE‑2026‑50751 (Active Bypass) THREAT: 5 Disable IKEv1 and legacy clients, require machine certificates, and patch. Hunt for suspicious RA logins since at least 2026‑05‑07. Remote Access IKEv1 JPCERT Alert C:ESCALATE PAN‑OS GlobalProtect CAS Authentication Bypass CVE‑2026‑0265 THREAT: 4 Remove CAS profiles from all login interfaces until patched; restrict exposure; enforce MFA or device certs; monitor anomalous VPN sessions. Palo Alto Networks GlobalProtect JPCERT Alert C:ESCALATE EVbee DC‑80 EV Chargers: Multi‑Vector RCE/Control THREAT: 5 Unauthenticated web/Bluetooth control, unsigned firmware, and OCPP injection enable root RCE and unsafe operations. Isolate chargers and enforce OCPP/mTLS hardening. OCPP Firmware Signing CISA Bulletin C:ESCALATE ICU Scandinavia Boomerang (Monitoring) Critical Flaws THREAT: 4 Plaintext service/SMTP creds via HTTP and unauthorized writes to sensor DBs risk data integrity and disruption. Upgrade ≥2.4.18.029 and rotate all credentials. Credential Exposure Config Leakage CERT Polska C:ESCALATE China‑Nexus Covert Networks (SOHO/IoT Proxies) THREAT: 4 Nation‑state actors route ops via compromised edge devices for stealthy CNI access. Harden routers/IoT, enable egress controls and NetFlow/NDR. Volt Typhoon Egress Filtering Joint Advisory C:ESCALATE FortiGate Credential Leak (“FortiBleed”) THREAT: 4 Credentials for >86k devices circulating; evidence of valid reuse. Rotate secrets, enforce PBKDF2, audit for config exports and AD abuse. JPCERT Alert C:ESCALATE Nation‑State Campaigns Targeting Routers/Edges THREAT: 5 PRC and Russian operations modify ACLs, enable unusual mgmt ports, abuse tunnels, and capture AAA creds on Cisco/others—risk of OT pivoting. Cisco IOS/IOS‑XE AAA Harvesting CSA (PRC) CSA (FSB) C:ESCALATE KEV Updates: SonicWall SMA1000, ADFS, SharePoint THREAT: 4 Actively exploited edge/identity CVEs; isolate, patch, validate pre‑patch compromise. CISA KEV C:ESCALATE KEV: FortiSandbox (Cmd Injection) & SharePoint THREAT: 4 Patch Fortinet appliances; restrict admin exposure; hunt for web‑to‑process activity on SharePoint. CISA KEV C:ESCALATE KEV: KNX (Building Automation) & Oracle EBS THREAT: 4 Segment KNX/IP, enable KNX Secure; apply Oracle CPU; hunt for privilege abuse. CISA KEV C:ESCALATE GeoServer RCE—Gov’t Incident Lessons & Active Botnet Abuse THREAT: 4 Exploited CVE‑2024‑36401 used for webshells and lateral movement; Rondo botnet activity observed in the wild. CISA Advisory SANS ISC C:ESCALATE ICS‑Relevant Weekly CVE Summaries THREAT: 4–5 Multiple OT/ICS CVEs (Delta PLC, OpENer, PcVue, DAQFactory). Prioritize patch planning and segmentation. SB26‑187 SB26‑194 SB26‑180 C:ESCALATE OilRig/APT34 “Project CAV3RN” (Graph Calendar C2) THREAT: 4 Abuse of Outlook calendar events and DNS AAAA for stealthy C2; rotate service principals and hunt for calendar anomalies. Kaspersky GReAT C:ESCALATE Fortinet/SharePoint: Additional KEV Items THREAT: 3–4 New KEV entries mandate rapid remediation and compromise checks per BOD 26‑04. CISA KEV C:ESCALATE Pro‑Russia Hacktivists Target Exposed VNC on OT/HMI THREAT: 4 Opportunistic setpoint changes and alarm suppression. Remove VNC exposure; require VPN+MFA for remote ops. Joint Advisory B:MONITOR B:MONITOR Siemens RUGGEDCOM APE1808 hosting PAN‑OS Issues THREAT: 3 PAN‑OS vulnerabilities affect virtual NGFW on APE; patch per Palo Alto advisories; restrict mgmt exposure. ICS Advisory B:MONITOR Siemens SIDIS Secured SmartPlug Critical Vulns THREAT: 4 Upgrade to v7.26.0310; segment and restrict SSH/TLS management; rotate keys post‑update. ICS Advisory B:MONITOR Tycon TPDIN‑Monitor‑WEB2 Auth Bypass/Cleartext Secrets THREAT: 3 Remove public web access, rotate device credentials, and apply vendor mitigations. ICS Advisory B:MONITOR TrickBot Variant Using DNS Tunneling C2 THREAT: 3 Inspect DNS egress, restrict public resolvers, and monitor scheduled task/ADS abuse tied to ransomware operators. View Article B:MONITOR Ubuntu snap‑confine LPE (CVE‑2026‑8933) THREAT: 3 Patch snapd packages; watch /run/udev rules and FUSE mounts; LPE can convert phishing footholds to root on jump hosts. View Article B:MONITOR FBI/IRS‑CI Seize NetNut; TV Proxy Botnet Disrupted THREAT: 3 Block residential proxy exits; monitor auth anomalies; proxy capacity may reconstitute via resellers. KrebsOnSecurity B:MONITOR Atlassian Stack: Multiple Critical/High Fixes THREAT: 3 Patch Bamboo/Bitbucket/Confluence/Jira; protect CI/CD tokens to prevent supply‑chain impact on OT gateways. View Article B:MONITOR HelloNet Campaign via ViPNet Update System THREAT: 4 DLL sideloading/injection to svchost; rotate Entra/app creds; segment ViPNet bridges to OT. Kaspersky B:MONITOR Kaspersky ICS Q1‑2026 Telemetry THREAT: 3 Elevated exposure for electric power, construction, building automation; emphasize email/web filtering and AutoCAD controls. ICS CERT Report A:IGNORE A:IGNORE NCSC Policy/AI Overviews (No OT Indicators) THREAT: 1 High‑level guidance without actionable IOCs or OT/ICS relevance. View Article View Article A:IGNORE Captive Portal Detection Traffic Explained THREAT: 1 Benign OS/browser connectivity checks; allowlist known endpoints to reduce false positives. SANS ISC Source References AA26‑097a: Iranian‑Affiliated Cyber Actors Exploit PLCs — Security Advisory SharePoint Hardening/Exploitation — CISA Alert Check Point RA CVE‑2026‑50751 — JPCERT PAN‑OS CAS CVE‑2026‑0265 — JPCERT EVbee DC‑80 — CISA Bulletin ICU Boomerang — CERT Polska China‑Nexus Covert Networks — Joint Advisory FortiBleed — JPCERT Router/Edge CSAs — PRC CSA | FSB CSA KEV Updates (7/14, 7/15, 7/16, 7/21, 7/22) — 7/14 7/15 7/16 7/21 7/22 ICS CVE Summaries — SB26‑187 SB26‑180 SB26‑194 Tycon / Siemens ICS Advisories — Tycon Siemens SIDIS Siemens APE1808 Siemens IAM Client GeoServer IR & SANS Botnet Activity — CISA SANS ISC APT34 CAV3RN (Kaspersky) — Securelist Generated from AI-assisted RSS OT/ICS cyber threat intelligence analysis. All source URLs preserved and linked in-line."
}
